# packaging/systemd/aegis.service
[Unit]
Description=Aegis - Secure Infrastructure DPU Agent
Documentation=https://github.com/gobeyondidentity/cobalt
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/secureinfra/aegis.env
ExecStart=/usr/bin/aegis --server $AEGIS_SERVER --dpu-name $AEGIS_DPU_NAME \
    ${DOCA_PCI_ADDR:+--doca-pci-addr $DOCA_PCI_ADDR} \
    ${DOCA_REP_PCI_ADDR:+--doca-rep-pci-addr $DOCA_REP_PCI_ADDR} \
    ${DOCA_SERVER_NAME:+--doca-server-name $DOCA_SERVER_NAME}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=LD_LIBRARY_PATH=/opt/mellanox/doca/lib/aarch64-linux-gnu

# DPU agent needs elevated privileges for hardware access
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target
