# packaging/systemd/aegis.service
[Unit]
Description=Aegis - Secure Infrastructure DPU Agent
Documentation=https://github.com/gobeyondidentity/secure-infra
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/secureinfra/aegis.env
ExecStart=/usr/bin/aegis \
    --server ${AEGIS_SERVER} \
    --dpu-name ${AEGIS_DPU_NAME} \
    --listen ${AEGIS_LISTEN:-:18051} \
    --local-listen ${AEGIS_LOCAL_LISTEN:-localhost:9443} \
    --db-path ${AEGIS_DB_PATH:-/var/lib/aegis/aegis.db} \
    --keystore ${AEGIS_KEYSTORE:-/var/lib/secureinfra/known_hosts.json}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=LD_LIBRARY_PATH=/opt/mellanox/doca/lib/aarch64-linux-gnu

# DPU agent needs elevated privileges for hardware access
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target
