# packaging/systemd/sentry.service
[Unit]
Description=Sentry - Secure Infrastructure Host Agent
Documentation=https://github.com/gobeyondidentity/cobalt
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=secureinfra
Group=secureinfra

# Load environment configuration (optional file)
EnvironmentFile=-/etc/default/sentry

# CLI flags (sentry does not support --config)
# Environment variables from /etc/default/sentry:
#   SENTRY_FORCE_COMCH=1         -> --force-comch
#   SENTRY_DOCA_PCI_ADDR=...     -> --doca-pci-addr=...
#   SENTRY_DOCA_SERVER_NAME=...  -> --doca-server-name=...
#   SENTRY_AUTH_KEY=...          -> --auth-key=...
ExecStart=/bin/sh -c '/usr/bin/sentry \
    ${SENTRY_FORCE_COMCH:+--force-comch} \
    ${SENTRY_DOCA_PCI_ADDR:+--doca-pci-addr=$SENTRY_DOCA_PCI_ADDR} \
    ${SENTRY_DOCA_SERVER_NAME:+--doca-server-name=$SENTRY_DOCA_SERVER_NAME} \
    ${SENTRY_AUTH_KEY:+--auth-key=$SENTRY_AUTH_KEY}'

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/secureinfra /etc/secureinfra

[Install]
WantedBy=multi-user.target
